# debian/postinst
#! /bin/bash
systemctl --global  enable --now sdgyrodsu.service || echo "sdgyrodsu service could not be enabled"
