# debian/prerm
#! /bin/bash
systemctl --global disable --now sdgyrodsu.service || echo "sdgyrodsu service could not be disabled"
